# config/ipmi-vmware-bridge.service
[Unit]
Description=IPMI VMware Bridge Service - OpenShift Virtualization Integration
Documentation=https://github.com/your-repo/ipmi-vmware-bridge
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/lchiaret/git/ipmi-vmware
ExecStart=/usr/bin/python3 /home/lchiaret/git/ipmi-vmware/src/ipmi_bridge.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

# Environment variables
Environment=IPMI_DEBUG=true
Environment=PYTHONPATH=/home/lchiaret/git/ipmi-vmware/src

# Security settings
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/log /var/run /home/lchiaret/git/ipmi-vmware
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ipmi-vmware-bridge
Environment=PYTHONPATH=/opt/ipmi-vmware-bridge
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
