# config/ipmi-vmware-bridge.service
[Unit]
Description=IPMI VMware Bridge Service - OpenShift Virtualization
Documentation=https://github.com/lchiaret/ipmi-vmware-bridge
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/lchiaret/git/ipmi-vmware
ExecStart=/usr/bin/python3 /home/lchiaret/git/ipmi-vmware/src/ipmi_bridge.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

# Security settings (relaxed for root execution)
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/var/log /var/run /home/lchiaret/git/ipmi-vmware

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ipmi-vmware-bridge

# Environment variables
Environment=PYTHONPATH=/home/lchiaret/git/ipmi-vmware/src
Environment=PYTHONUNBUFFERED=1
Environment=IPMI_DEBUG=false

[Install]
WantedBy=multi-user.target
