# config/redfish-vmware-server.service
[Unit]
Description=Redfish VMware Server - OpenShift Virtualization
Documentation=https://github.com/lchiaret/redfish-vmware-server
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=__PROJECT_ROOT__
Environment=PYTHONUNBUFFERED=1
Environment=REDFISH_DEBUG=true
Environment=REDFISH_PERF_DEBUG=true
Environment=REDFISH_VMWARE_DEBUG=true
Environment=PYTHONPATH=__PROJECT_ROOT__/src
ExecStart=/usr/bin/python3 __PROJECT_ROOT__/src/redfish_server.py --config __PROJECT_ROOT__/config/config.json
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

# Security settings (relaxed for root execution)
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/var/log /var/run __PROJECT_ROOT__

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redfish-vmware-server

[Install]
WantedBy=multi-user.target
